# hotkeys.service
[Install]
WantedBy=multi-user.target

[Unit]
After=display-manager.service

[Service]
Type=simple
ExecStart=/usr/bin/hotkeys.sh
